# launchd/ai.hermes-staging.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.smartclaw-staging</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>@REPO_ROOT@/scripts/launchd-env-wrapper.sh</string>
        <string>@HERMES_BIN@</string>
        <string>gateway</string>
        <string>run</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>@HOME@</string>
        <key>PATH</key>
        <string>@HERMES_EXTRA_PATH@@HOME@/.local/bin:@HOME@/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HERMES_HOME</key>
        <string>@HOME@/.smartclaw</string>
        <key>HERMES_LOG_LEVEL</key>
        <string>INFO</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>GATEWAY_ALLOW_ALL_USERS</key>
        <string>true</string>
        <key>SLACK_ALLOW_BOTS</key>
        <string>mentions</string>
        <key>SLACK_REQUIRE_MENTION</key>
        <string>true</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>@HOME@/.smartclaw</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>@HOME@/.smartclaw/logs/gateway.log</string>

    <key>StandardErrorPath</key>
    <string>@HOME@/.smartclaw/logs/gateway.error.log</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
